# client/supervisor/sensors.conf
[program:tunnel]
command=/home/pi/sensors.sh
user=pi
directory=/home/pi
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/ssh_tunnel.err
stdout_logfile=/var/log/ssh_tunnel.log
